# roles/gitian/templates/rc.local
#!/bin/sh -e

echo 1 > /sys/fs/cgroup/cpuset/cgroup.clone_children
echo 1 > /proc/sys/kernel/unprivileged_userns_clone
echo 1 > /proc/sys/net/ipv4/ip_forward

exit 0
